# .config/systemd/user/start-up-routine.service
[Unit]
Description=Start-up Routine

[Service]
Type=simple
Environment="SCRIPTS_PATH=$HOME/.scripts"
Environment="LOG_PATH=$HOME/.scripts.log"
ExecStart=/home/eugene/.scripts/on-start-up.sh

[Install]
WantedBy=default.target
