# .config/systemd/user/start-up-routine.service
[Unit]
Description=Start-up Routine

[Service]
Type=simple
Environment="SCRIPTS_PATH=/home/eugene/.scripts"
Environment="LOG_PATH=/home/eugene/.scripts.log"
ExecStart=/bin/bash /home/eugene/.scripts/on-start-up.sh

[Install]
WantedBy=default.target
