# .config/systemd/user/start-up-routine.service
[Unit]
Description=Start-up Routine

[Service]
Type=simple
ExecStart=/home/eugene/.scripts/on-start-up.sh

[Install]
WantedBy=default.target
