# .config/systemd/user/wezterm-tabs-cache.service
[Unit]
Description=Update WezTerm tabs cache

[Service]
Type=oneshot
ExecStart=%h/.scripts/wezterm-tabs-update.sh
