# .scripts/udev-rules/10-monitor.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/eugene/.scripts/monitor-setup-hook.sh"
